# deploy/jarvis.service
[Unit]
Description=Jarvis-AI Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/Jarvis-AI
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/Jarvis-AI/.venv/bin/python launcher.py
Restart=always
RestartSec=5

# Hardening (optional)
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
